# build/rhem-server/etc/containers/systemd/postgresql.container
[Unit]
Description=PostgreSQL for Keycloak

[Container]
Label=app=keycloak
Image=quay.io/sclorg/postgresql-16-c9s:latest
ContainerName=keycloak-db
Network=flightctl.network
Environment=POSTGRES_DB=keycloak
Environment=POSTGRES_USER=keycloak
Environment=POSTGRES_PASSWORD=keycloak

Volume=pg-data:/var/lib/postgresql/data:Z

PublishPort=5432:5432

[Install]
WantedBy=sso.target
