# build/rhem-server/usr/lib/systemd/system/sso.target
[Unit]
Description=FlightCTL SSO Group

# Requires= directives ensures these services are started with this unit.
# After= directives delays this unit until these services start.
#
# Note: The exact startup order among the services is controlled by
# additional directives defined within the .container files
Requires=flightctl-network.service

After=network.target
After=flightctl-network.service

[Install]
WantedBy=multi-user.target